# setup/emu904/scripts/EmuServiceGuard.cron
# Check for dead xdaqd services and resurrect them.
*/1 * * * * root /opt/xdaq/share/`cat /etc/xdaq.zone`/scripts/EmuServiceGuard.sh
